# tags/r3-6-4/examples/logrotate.cron
#!/bin/sh

/usr/sbin/logrotate /etc/logrotate.conf
